# salt/salt-masterless.service
[Unit]
Description=Salt masterless periodic salt-call

[Service]
Type=oneshot
ExecStart=/usr/bin/salt-call state.apply